# deploy/launchd/com.user.mem0-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.user.mem0-watchdog</string>

	<key>ProgramArguments</key>
	<array>
		<string>__HELIX_DEV_TOOLS_BIN__</string>
		<string>mem0-watchdog</string>
		<string>--probe-url</string>
		<string>http://127.0.0.1:${MEM0_PORT:-8888}/docs</string>
		<string>--interval</string>
		<string>60s</string>
		<string>--timeout</string>
		<string>5s</string>
		<string>--threshold</string>
		<string>3</string>
		<string>--log-path</string>
		<string>__HOME__/logs/runx/mem0-watchdog.ndjson</string>
	</array>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<true/>

	<key>ThrottleInterval</key>
	<integer>30</integer>

	<key>StandardOutPath</key>
	<string>__HOME__/logs/runx/mem0-watchdog.stdout.log</string>

	<key>StandardErrorPath</key>
	<string>__HOME__/logs/runx/mem0-watchdog.stderr.log</string>

	<key>ProcessType</key>
	<string>Background</string>
</dict>
</plist>
